# p4mgr.service
[Unit]
Description=P4Mgr LED Matrix Display Manager
After=network.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/p4mgr
ExecStart=/usr/bin/python3 /home/pi/p4mgr/main.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target