# p4mgr.service.example
[Unit]
Description=P4Mgr LED Matrix Display Manager
After=network.target

[Service]
Type=simple
User=mizphses
Group=mizphses
WorkingDirectory=/home/mizphses/p42
ExecStart=/usr/bin/python3 /home/mizphses/p42/main.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target